# plexpass/services.ini
[program:plexmediaserver]
user = user
autorestart = true
command = /start_pms
