# sonarr/services.ini
[program:sonarr]
user = user
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
command = sonarr -nobrowser -data=/config
