# sonarr/services.ini
[program:sonarr]
user = user
autorestart = true
command = /usr/bin/sonarr -nobrowser -data=/config
