# headphones/services.ini
[program:headphones]
user = user
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
command = /usr/bin/python2 /opt/headphones/Headphones.py --nolaunch --datadir /config
