# subliminal-auto/services.ini
[program:subliminal_watchdog]
user = user
autorestart = true
command = /subliminal_watchdog.py
