# couchpotato/services.ini
[program:couchpotato]
user = user
autorestart = true
command = /CouchPotatoServer/CouchPotato.py --data_dir /config/data --config_file=/config/CouchPotato.cfg --console_log
