# couchpotato/services.ini
[program:couchpotato]
user = user
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
command = /CouchPotatoServer/CouchPotato.py --data_dir /config/data --config_file=/config/CouchPotato.cfg --console_log
