# nzbget/services.ini
[program:nzbget]
user = user
autorestart = true
command = /opt/nzbget/nzbget -c /config/nzbget.conf -o OutputMode=log -s
